# backend/deploy/bmtc-backup.service
[Unit]
Description=BMTC API Database Backup

[Service]
Type=oneshot
User=bmtc
Group=bmtc
EnvironmentFile=/etc/bmtc-api/env
ExecStart=/opt/bmtc-api/scripts/backup.sh
