# backend/deploy/bmtc-retention.timer
[Unit]
Description=BMTC API Daily Retention Cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
